# deployments/swoop-dr-paper/service-files/systemd/claude-orchestrator.service
[Unit]
Description=Dr Paper Orchestrator (claude-code-by-agents)
After=network-online.target claude-agentapi.service tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=claude-agent
WorkingDirectory=/opt/claude-code-by-agents/backend
EnvironmentFile=/etc/claude/api-key
Environment="DENO_DIR=/home/claude-agent/.deno"
Environment="PATH=/home/claude-agent/.deno/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/claude-agent/.deno/bin/deno task dev
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dr-paper

[Install]
WantedBy=multi-user.target
