# deployments/swoop-dr-paper/service-files/systemd/claude-agentapi.service
[Unit]
Description=Claude Code AgentAPI (Dr Paper fleet)
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=claude-agent
WorkingDirectory=/home/claude-agent
EnvironmentFile=/etc/claude/api-key
ExecStart=/usr/local/bin/agentapi server --type claude -- claude \
  --dangerously-skip-permissions \
  --max-budget-usd 10.00
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-agentapi

[Install]
WantedBy=multi-user.target
